# projects/ROCKNIX/devices/H700/filesystem/usr/lib/udev/rules.d/wlan0.powersave.off.rules
ACTION=="add", SUBSYSTEM=="net", SUBSYSTEMS=="sdio", KERNELS=="mmc2:0001:1", RUN+="/sbin/iw dev %k set power_save off"
